# root/etc/systemd/system/dnsmasq.service.d/ipset.conf
[Service]
ExecStartPre=/etc/e-smith/events/actions/nethserver-squid-ipset
